# brew/brew-base/etc/profile.d/ansible.sh
export ANSIBLE_CONFIG="/Users/potetm/.ansible.cfg"

ANSIBLE_HOME="/Users/potetm/projects/ansible/ansible"
. "$ANSIBLE_HOME"/hacking/env-setup -q > /dev/null
