# brew/brew-base/etc/profile.d/brew-sbin.sh
SBIN_BIN="/usr/local/sbin"
PATH="$SBIN_BIN":"$PATH"

export PATH
