# brew/brew-base/etc/profile.d/elasticsearch.sh
PATH="/opt/elasticsearch/home/bin:$PATH"
export PATH
