# brew/brew-base/etc/profile.d/gcloud.sh
PATH="/opt/gcloud/home/bin:$PATH"
export PATH
