# brew/brew-base/etc/profile.d/go.sh
GOPATH="$HOME/go"
PATH="$HOME/go/bin:$PATH"
export GOPATH
export PATH
