# brew/brew-base/etc/profile.d/java.sh
export JAVA_HOME="/opt/java/home"
export JVM_OPTS="-Xmx6g"
