# brew/brew-base/etc/profile.d/java.sh
export JAVA_HOME="/opt/java/home"
export PATH="$JAVA_HOME/bin:$PATH"
