# brew/brew-base/etc/profile.d/pyenv-virtualenv.sh
# Because pyenv-virtualenv adds to $PROMPT_COMMAND, and I don't like it.

#if which pyenv-virtualenv-init > /dev/null; then
#  eval "`pyenv virtualenv-init -`";
#fi
