# brew/brew-base/etc/profile.d/rabbitmq.sh
RABBITMQ_HOME="/usr/local/opt/rabbitmq"
PATH="$RABBITMQ_HOME/sbin:$PATH"

export RABBITMQ_HOME
export PATH
