# brew/brew-base/etc/profile.d/scala.sh
export SCALA_HOME="/opt/scala/home"
export PATH="$SCALA_HOME/bin:$PATH"
