# brew/brew-base/etc/profile.d/userbin.sh
USER_BIN="`echo ~/bin`"
USER_BIN_NO_LINK="`echo ~/bin-no-conf-link`"
PATH="$USER_BIN":"$USER_BIN_NO_LINK":"$PATH"

export PATH
